# conf/crontab
0 0 * * * echo "generate:all" > "${CAVEDB_WORKER_FIFO}"
0 1 * * * echo "backup" > "${CAVEDB_WORKER_FIFO}"
