# conf/crontab
0 0 * * * root touch "__CAVEDB_WORKER_MSG_DIR__/elevation_dem_update"

0 2 * * * root touch "__CAVEDB_WORKER_MSG_DIR__/backup"

