# systemd/build-all-bulletins.timer
[Unit]
Description=Build all bulletins

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
